# covid-reporter.service
[Unit]
Description=Covid reporter

[Service]
Type=simple
ExecStart=/home/pi/ti-covid/start.sh

[Install]
WantedBy=multi-user.target
